# clean-rice/.config/systemd/user/hyprshade.service
[Unit]
Description=Apply hyprshade auto configuration

[Service]
Type=oneshot
ExecStart=/usr/bin/hyprshade auto

[Install]
WantedBy=default.target
